# systemd/dev_maintenance.service
[Service]
ExecStart=/home/vmagro/bin/dev_maintenance
Type=oneshot
