# systemd/dev_maintenance.timer
[Unit]
Description=Run dev_maintenance every day at 4am

[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true

[Install]
WantedBy=timers.target
