# systemd/dev_maintenance.service
[Service]
ExecStart=/home/vmagro/bin/dev_maintenance configerator www opsfiles fbcode
Type=oneshot
OnFailure=status-email-vmagro@%n.service
