# systemd/xscreensaver.service
[Unit]
Description=XScreenSaver
[Service]
ExecStart=/usr/bin/xscreensaver
Restart=Always
[Install]
WantedBy=default.target
